# docker/conf/supervisor/caddy-supervisor.conf
[program:nginx]
command=/usr/bin/caddy
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
priority=10
user=shopware
process_name=shopware-%(program_name)s
