# docker/conf/supervisor.conf
[supervisord]
nodaemon=true
logfile=/dev/stderr
logfile_maxbytes=0

[program:php-fpm]
command=/usr/local/sbin/php-fpm
redirect_stderr=true
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /usr/local/etc/nginx/nginx.conf -g "daemon off;"
redirect_stderr=true
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0